# system_files/usr/lib/systemd/user/hyprland-shutdown.target
[Unit]
Description=Shutdown running hyprland session
DefaultDependencies=no
StopWhenUnneeded=true

Conflicts=graphical-session.target graphical-session-pre.target
After=graphical-session.target graphical-session-pre.target
OnFailureJobMode=replace-irreversibly
